# ingest.cron
iPATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
0 15 * * * cd /code && /usr/local/bin/python /usr/local/bin/flask --app=dgrehydro ingest_riverine > /proc/1/fd/1 2>/proc/1/fd/2
0 */6 * * * cd /code && /usr/local/bin/python /usr/local/bin/flask --app=dgrehydro ingest_flashflood > /proc/1/fd/1 2>/proc/1/fd/2
*/10 * * * * cd /code && /usr/local/bin/python /usr/local/bin/flask --app=dgrehydro ingest_dustwarning > /proc/1/fd/1 2>/proc/1/fd/2
